# scripts/com.zor125.paperagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zor125.paperagent</string>

  <key>RunAtLoad</key>
  <true/>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/zor125/Projects/Paper_Agent/scripts/run_daily.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/zor125/Projects/Paper_Agent</string>

  <key>StandardOutPath</key>
  <string>/Users/zor125/Projects/Paper_Agent/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/zor125/Projects/Paper_Agent/logs/launchd.err.log</string>
</dict>
</plist>
